# systemd/month-dir.timer
[Unit]
Description=Create current month dir on ~/tmp

[Timer]
Persistent=true
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
